# DEFAULT-CONFIGS/etc/profile.d/gtk3-noosb.sh
#/bin/sh
## Disables the overly scrollbars in GTK3 applications for the current user
export GTK_OVERLAY_SCROLLING=0
